# isecl-k8s-controller/dist/linux/isecl-k8s-controller.service
[Unit]
Description=ISecL k8s custom controller
ConditionPathExists=/opt/isecl-k8s-extensions/isecl-k8s-controller/bin/isecl-k8s-controller

[Service]
ExecStart=/usr/bin/isecl-k8s-controller -kubeconf=/etc/kubernetes/admin.conf -log_dir=/var/log/messages
KillMode=process
Restart=on-failure
RestartPreventExitStatus=255
Type=simple

[Install]
WantedBy=multi-user.target
Alias=isecl-k8s-controller.service

